# files/oci/tuning/dbus-broker.service.d/10-cosmic-nofile.conf
# Raise dbus-broker's file-descriptor limit. The systemd default of
# 1024 lands COSMIC in trouble: cosmic-comp + cosmic-settings-daemon +
# cosmic-applets + every desktop application all sit on dbus, and a
# few hundred file watchers across cosmic-files / cosmic-edit windows
# are enough to push the broker past the soft limit and start
# rejecting connections. Pop!_OS ships this exact drop-in for the
# same reason on COSMIC.
[Service]
LimitNOFILE=524288:1048576
