# files/oci/tuning/dbus.service.d/10-cosmic-nofile.conf
# Raise the system D-Bus daemon's file-descriptor limit. The systemd
# default of 1024 lands COSMIC in trouble: cosmic-comp + cosmic-
# settings-daemon + cosmic-applets + every desktop application all
# sit on the system bus, and a few hundred file watchers across
# cosmic-files / cosmic-edit windows are enough to push the daemon
# past the soft limit and start rejecting connections. Pop!_OS ships
# the same drop-in -- theirs targets dbus-broker (Fedora default);
# fdsdk packages the reference dbus-daemon, so we put it here.
[Service]
LimitNOFILE=524288:1048576
